# s6-overlay/s6-rc.d/openldap/run
#!/command/with-contenv sh
# slapd -h $SLAPD_HOSTS -F $SLAPD_CONF_DIR -u ldap -g ldap -d $CONTAINER_LOG_LEVEL >> /dev/stdout 2>&1
while true; do echo "test"; sleep 200000; done
